# files/etc/supervisor/conf.d/logs.conf
[program:logs]
command=/usr/bin/tail -F /var/log/nginx/access.log /var/log/nginx/error.log 
autostart=true
autorestart=true
startretries=3
priority=5
